# ansible/roles/rpi_mqtt_services/templates/mqtt-led.service.j2
[Unit]
Description=MQTT LED Service - GPIO Control
After=mosquitto.service
Requires=mosquitto.service

[Service]
Type=notify
User={{ ansible_user }}
WorkingDirectory={{ mqtt_led_dir }}
Environment=GPIOZERO_PIN_FACTORY=lgpio
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 {{ mqtt_led_dir }}/mqtt_led.py
WatchdogSec=60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
